# apps/cms/cmd/cli/nien-su-viet-cms.service
[Unit]
Description=Nien Su Viet CMS Service
After=network.target postgresql.service rabbitmq-server.service
Wants=postgresql.service rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/nien-su-viet-cms
ExecStart=/opt/nien-su-viet-cms/server
Restart=always
RestartSec=5
StandardOutput=append:/var/log/nien-su-viet-cms/app.log
StandardError=append:/var/log/nien-su-viet-cms/error.log

# Environment variables (override with systemd drop-in file)
Environment="GIN_MODE=release"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/nien-su-viet-cms
ReadOnlyPaths=/opt/nien-su-viet-cms

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
